# uinput.service
[Unit]
Description=Fix uinput permissions

[Service]
ExecStart=chmod +0666 /dev/uinput

[Install]
WantedBy=default.target
